# 80-koy.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/systemctl start koy@%k.service"
