# 80-kou.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/systemctl start kou@%k.service"
